# systemd/.config/systemd/user/startup-script.service
[Unit]
Description=Run my startup script

[Service]
ExecStart=/home/edret/scripts/run.sh

[Install]
WantedBy=default.target

